# mcp-server/com.fluxa-map.mcp.local.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.fluxa-map.mcp.local</string>

  <key>ProgramArguments</key>
  <array>
    <string>/tmp/run-fluxa-map-mcp.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/williamwang/Downloads/项目/Fluxa Map/mcp-server</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/fluxa-map-mcp-launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/tmp/fluxa-map-mcp-launchd.err.log</string>
</dict>
</plist>
